# systemd/maintenance.service
[Unit]
Description=Maintenance Agent
Documentation=https://github.com/swordh/Maintenance
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/maintenance/agent/maintenance_agent.py
WorkingDirectory=/opt/maintenance/agent
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target
